# scripts/systemd/ttserver-watcher.path
[Path]
Unit=ttserver-watcher.service
PathChanged=/home/ttserver/dist/tt-server

[Install]
WantedBy=multi-user.target

